# website/config/supervisor/apache2.conf
[program:apache2]
command=service apache2 start
user=root
autostart=true
stdout_logfile=/var/log/supervisor/apache2.log
redirect_stderr=true
